# remind.service
[Unit]
Description=WhatsApp Reminder Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/whatsapp-bot
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
